# supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=debug
logfile_maxbytes=1000000
logfile_backups=1

; db-init runs first (priority 10) — applies migrations on first boot, then exits
[program:db-init]
command=bash /app/db-init.sh
autorestart=false
startsecs=0
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; PostgREST starts after db-init (priority 20)
[program:postgrest]
command=bash /app/postgrest-entrypoint.sh
autorestart=unexpected
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; n8n starts after PostgREST (priority 30)
[program:n8n]
command=bash /app/n8n-entrypoint.sh
autorestart=unexpected
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; nginx starts last among core services (priority 40)
[program:nginx]
command=bash /app/nginx-entrypoint.sh
autorestart=unexpected
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; workflow-import runs once after n8n is up (priority 50)
; It polls /healthz internally so it's safe to start early
[program:workflow-import]
command=bash /app/workflow-import.sh
autorestart=false
startsecs=0
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ha-integration runs once after workflow-import (priority 60)
; Auto-configures Home Assistant to call the n8n-claw agent
; Requires hassio_api: true + homeassistant_api: true in config.json
; Skipped gracefully if SUPERVISOR_TOKEN is not set (local dev mode)
[program:ha-integration]
command=bash /app/ha-integration.sh
autorestart=false
startsecs=0
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
